# dist/systemd/evo.service.d/exec-start.conf
# Systemd unit drop-in: override the framework reference
# unit's `ExecStart=/opt/evo/bin/evo` so the steward launches
# this distribution's reference binary `evo-device-audio`
# instead.
#
# The framework reference unit bakes
# `ExecStart=/opt/evo/bin/evo` because the framework's own
# steward binary is named `evo`. Each domain distribution
# ships a differently-named steward binary that links the
# framework crate plus the distribution's plugin set; for
# this distribution that binary is `evo-device-audio`.
#
# Drop-in mechanics: systemd treats `ExecStart=` as a
# list-typed directive. The empty `ExecStart=` line clears
# any inherited entry from the base unit, and the second
# `ExecStart=` line sets the effective command. Without the
# empty-clear, systemd would append rather than replace.
#
# Apply: drop this file at
# `/etc/systemd/system/evo.service.d/`, run `systemctl
# daemon-reload && systemctl restart evo.service`.
#
# Vendor distributions that ship their own steward binary
# under a different filename adapt this drop-in by editing
# the path. Distributions that prefer to ship a complete
# unit file rather than layer a drop-in can omit this file
# and place a full `evo.service` at
# `/etc/systemd/system/evo.service`.

[Service]
ExecStart=
ExecStart=/opt/evo/bin/evo-device-audio
