# dist/systemd/evo.service.d/https.conf
# Systemd unit drop-in: HTTPS substrate for this
# distribution. The framework's reference unit declares no
# HTTPS environment because the framework binary is
# transport-agnostic; each distribution that wants an HTTPS
# listener (alongside or instead of the unix socket) ships
# the env-var trio below.
#
# Why this lives at distribution scope and not in the
# framework's reference unit: a distribution that targets
# embedded MCUs may run the steward without HTTPS at all
# (unix socket only); a distribution that targets server
# tier may bind a different port or split static-asset
# serving onto a separate origin. The framework refuses to
# bake a transport choice into its reference unit.
#
# Audio reference picks 0.0.0.0:8443 (canonical HTTPS, all
# interfaces) so the operator's browser on the LAN reaches
# the device at `https://<hostname or ip>:8443/`. The state
# directory `/var/lib/evo/https` is StateDirectory= adjacent
# (the framework writes its CA, signing keys, and
# bootstrap-token material under
# `/var/lib/evo/https/https/`). The static-asset directory
# `/opt/evo/ui` is where the canonical UI install script
# writes the shell distribution (manifest + service worker
# + static bundle) so the same origin serves both the
# framework API and the operator UI shell.
#
# Drop this file at `/etc/systemd/system/evo.service.d/`,
# run `systemctl daemon-reload && systemctl restart
# evo.service`. Apply alongside `exec-start.conf`. The
# bootstrap.sh script installs this automatically when run
# with the default `EVO_INSTALL_SYSTEMD_DROP_INS=1`.
#
# Vendor distributions adjust the listen socket
# (e.g. `127.0.0.1:8443` for proxy-fronted deployments) or
# omit this drop-in entirely (unix-socket-only deployment)
# by editing the values below or by overriding via a
# subsequent drop-in.

[Service]
Environment="EVO_HTTPS_LISTEN_ADDR=0.0.0.0:8443"
Environment="EVO_HTTPS_STATE_DIR=/var/lib/evo/https"
Environment="EVO_HTTPS_STATIC_DIR=/opt/evo/ui"
