# dist/systemd/evo.service.d/mpd-restart-privileges.conf
# Systemd unit drop-in: relax the framework's reference
# hardening enough for the playback.mpd plugin's
# privileged-action surface to function on this
# distribution.
#
# The framework's reference `evo.service.example` bakes in
# `ProtectSystem=strict` + `NoNewPrivileges=yes` ŌĆö
# appropriate for a steward whose plugin set does not touch
# the host's privileged surface. This distribution's plugin
# set includes `org.evoframework.playback.mpd`, whose
# fragment-writer worker must:
#
#   1. write the framework-controlled MPD `audio_output`
#      block to `/etc/evo/mpd.conf` (so MPD's
#      `include_optional /etc/evo/mpd.conf` from
#      `/etc/mpd.conf` picks it up), and
#   2. restart `mpd.service` via
#      `sudo -n /usr/bin/systemctl restart mpd` after every
#      rewrite (so MPD re-reads the new device + format).
#
# Each of those needs one specific framework-hardening
# relaxation, applied here at distribution scope:
#
#   - `ReadWritePaths=/etc/evo` ŌĆö narrowest possible write
#     widening; only `/etc/evo` becomes writable to the
#     steward, not the rest of `/etc`. The fragment is the
#     single file the plugin writes; future framework-managed
#     configuration falling under `/etc/evo` follows the
#     same posture.
#   - `NoNewPrivileges=no` ŌĆö required for sudo's setuid bit
#     to elevate. The narrow sudoers drop-in shipped at
#     `/etc/sudoers.d/evo-mpd-restart` allows only one exact
#     command (`/usr/bin/systemctl restart mpd`) for the
#     steward service user. No `NOPASSWD: ALL`; no
#     interactive auth in the service path.
#   - `Environment=EVO_SYSTEMCTL=/usr/bin/systemctl` ŌĆö
#     declares the systemctl binary path the plugin invokes.
#     The same path appears in the sudoers drop-in; the
#     plugin reads this env var via `EVO_SYSTEMCTL` and uses
#     it for sudo's argv. Operators on a distribution with a
#     non-standard systemctl path (e.g. Yocto custom prefix)
#     edit both files in lockstep.
#
# Vendor distributions that prefer to ship the framework's
# unrelaxed hardening posture can drop this conf and either
# (a) run the steward as root (no sudo needed, restart goes
# direct) or (b) park the playback.mpd plugin's restart leg
# (worker publishes Failed status; previous fragment file
# stays in place; MPD picks up the fragment on its next
# manual restart).
#
# Apply: drop this file at
# `/etc/systemd/system/evo.service.d/`, run `systemctl
# daemon-reload && systemctl restart evo.service`.

[Service]
NoNewPrivileges=no
ReadWritePaths=/etc/evo
Environment=EVO_SYSTEMCTL=/usr/bin/systemctl
