# dist/systemd/evo.service.d/state-dir-mode.conf
# Systemd unit drop-in: widen StateDirectoryMode to 0755 so
# system services co-resident with the steward (typically
# MPD via this distribution's playback.mpd warden) can
# traverse `/var/lib/evo/` to reach plugin-managed
# subdirectories underneath. The framework reference unit
# defaults to 0750, which is correct for the framework's
# own state but blocks an external service like the
# `mpd` system user from reading the music tree at
# `/var/lib/evo/music/INTERNAL/...`.
#
# Drop this file at `/etc/systemd/system/evo.service.d/`
# alongside the framework's reference unit, then run
# `systemctl daemon-reload && systemctl restart evo.service`
# to apply. Persists across boots; idempotent.
#
# Vendor distributions that prefer to keep the music tree
# in its own state directory (e.g. `/var/lib/evo-music/`)
# can omit this drop-in and route MPD's music_directory
# at the alternate path; in that case the framework's
# default 0750 is sufficient.

[Service]
StateDirectoryMode=0755
